# docker/etc/services.d/anisync/run
#!/usr/bin/with-contenv bash

cd /opt/anisync || exit

exec \
	s6-setuidgid abc /opt/anisync/AniSync --DataDirectory "/config"